# fastapi.service
[Unit]
Description=Traditional China Tours - FastAPI Application
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/traditionalchina
Environment="PATH=/var/www/traditionalchina/venv/bin"
ExecStart=/var/www/traditionalchina/venv/bin/uvicorn main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
